# deploy/ubuntu/systemd/cmsfleet-worker.service
[Unit]
Description=cmsfleet Integration Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=buscms
Group=buscms
WorkingDirectory=/opt/cmsfleet/current
EnvironmentFile=-/etc/cmsfleet/integration-worker.env
ExecStart=/usr/bin/npm --workspace @cmsfleet/integration-worker run start:prod
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/cmsfleet/current/config /var/log/cmsfleet /var/backups/cmsfleet
StandardOutput=append:/var/log/cmsfleet/integration-worker.log
StandardError=append:/var/log/cmsfleet/integration-worker-error.log
SyslogIdentifier=cmsfleet-worker

[Install]
WantedBy=multi-user.target