# deploy/ubuntu/systemd/cmsfleet-backup.timer
[Unit]
Description=Run cmsfleet backups nightly

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
Unit=cmsfleet-backup.service

[Install]
WantedBy=timers.target